# launchd/com.ahsgr.brief.eod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ahsgr.brief.eod</string>

    <key>ProgramArguments</key>
    <array>
        <string>VENV_PY</string>
        <string>REPO_ROOT/brief.py</string>
        <string>eod</string>
    </array>

    <!-- 17:00 local, weekdays (Weekday 1=Mon … 5=Fri) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>WorkingDirectory</key>
    <string>REPO_ROOT</string>

    <key>StandardOutPath</key>
    <string>HOME/Library/Logs/ahsgr-brief/brief-eod.log</string>
    <key>StandardErrorPath</key>
    <string>HOME/Library/Logs/ahsgr-brief/brief-eod.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>HOME</string>
    </dict>
</dict>
</plist>
